# projects/blinkled/myblinkled/files/myblinkled.service
[Unit]
Description=Blink LED
After=network.target

[Service]
ExecStart=/usr/bin/myblinkled
Restart=always

[Install]
WantedBy=default.target
